# helper/com.clashfox.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clashfox.helper</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/com.clashfox.helper</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardErrorPath</key>
    <string>/var/log/clashfox-helper.log</string>
    
    <key>StandardOutPath</key>
    <string>/var/log/clashfox-helper.log</string>
</dict>
</plist>
